# etc/update-motd.d/05-firmware
#!/bin/sh

printf "\n"
printf " Firmware: DKSL-90_2.1-0 \n"
