# etc/update-motd.d/05-firmware
#!/bin/sh

printf "\n"
printf " Firmware: `/etc/dksl-90-vars firmware` \n"
